# infra/systemd/direwolf-tail.service
# Direwolf CSV tailer pushing packets into the display backend
[Unit]
Description=Direwolf CSV Tailer
After=network-online.target direwolf.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/direwolf_display
Environment=UV_PROJECT_ENV=/opt/direwolf_display/.venv
ExecStart=/usr/bin/env uv run python tools/direwolf_csv_tail.py \
    --log-dir /var/log/direwolf \
    --post-url http://127.0.0.1:9090/v1 \
    --since-start
Restart=always
RestartSec=5
User=direwolf
Group=direwolf
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
