# infra/systemd/direwolf.service
[Unit]
Description=Direwolf APRS TNC with CSV logging
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=direwolf
Group=direwolf
WorkingDirectory=/opt/direwolf-display
EnvironmentFile=-/etc/direwolf/direwolf.env
ExecStart=/opt/direwolf-display/scripts/run_direwolf.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
