# infra/systemd/direwolf-display.service
# Direwolf Display map backend
[Unit]
Description=Direwolf Display Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/direwolf-display
Environment=UV_PROJECT_ENV=/opt/direwolf-display/.venv
ExecStart=/usr/bin/env uv run uvicorn app.main:create_app --factory --host 0.0.0.0 --port 9090
Restart=on-failure
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
# Adjust the user/group to an unprivileged account that can access Direwolf logs.
User=direwolf
Group=direwolf

[Install]
WantedBy=multi-user.target
